# .config/systemd/user/battery-check.service
[Unit]
Description=Check battery and warn if low

[Service]
Type=oneshot
ExecStart=/bin/bash /home/james/bin/battery-check
