# .config/systemd/user/iplayer-radio-pvr.timer
[Unit]
Description=Schedule the iPlayer radio PVR service

[Timer]
OnCalendar=*-*-* 11:00:00

[Install]
WantedBy=timers.target
