# .config/systemd/user/iplayer-pvr-lock-check.timer
[Unit]
Description=Schedule the iPlayer PVR lock check service

[Timer]
OnCalendar=*-*-* 12:00:00

[Install]
WantedBy=timers.target
