# .config/systemd/user/iplayer-tv-pvr.service
[Unit]
Description=Run iPlayer TV PVR fetch

[Service]
Type=oneshot
ExecStart=/bin/bash /home/james/services/scripts/iplayer-tv-pvr.sh
