# .config/systemd/user/poweroff.service
[Unit]
Description=Shutdown every day at specified time

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
